# .config/systemd/user/common-voice.timer
[Unit]
Description=Remind me to contribute to Common Voice daily

[Timer]
OnCalendar=23:15
Persistent=no

[Install]
WantedBy=timers.target

