# .config/systemd/user/toolsdb-backup@.timer
[Unit]
Description=Back up Toolforge database %I daily

[Timer]
OnCalendar=*-*-* 22:30:00

[Install]
WantedBy=timers.target
