# .config/systemd/user/toolsdb-backup@.service
[Unit]
Description=Back up Toolforge database %I

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-toolsdb %I
